# .cd/web-calc/etc/systemd/system/web-calc.service
[Unit]
Description=Web Calc
After=network.target
[Service]
Type=simple
ExecStart=/opt/goapp/build_go
[Install]
WantedBy=multi-user.target
