# .ci/web-calc/DEBIAN/postrm
#!/bin/bash
systemctl stop web-calc
systemctl disable web-calc
systemctl daemon-reload
echo alles

